# pump-control.service
[Unit]
Description=Pump Relay Control Web Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/pi-control-program
Environment=RELAY_PINS=27,22,23,24
Environment=RELAY_ACTIVE_LOW=1
Environment=GPIOZERO_PIN_FACTORY=lgpio
ExecStart=/home/pi/pi-control-program/.venv/bin/uvicorn app:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
